# traccar/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Traccar
# Runs the Nginx daemon
# ==============================================================================
bashio::net.wait_for 18682 localhost 900
bashio::log.info "Starting NGinx server..."
exec nginx
