# traccar/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Traccar
# Runs the Nginx daemon
# ==============================================================================
bashio::net.wait_for 18682
bashio::log.info "Starting NGinx server..."
exec nginx
