# traccar/rootfs/etc/services.d/traccar/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Traccar
# Runs the Traccar daemon
# ==============================================================================
declare -a options

bashio::log.info "Starting Traccar server..."

# Memory settings
options+=(-Xms256m)
options+=(-Xmx512m)

# Prefer IPv4, because of Hassio general issues with IPv6
options+=("-Djava.net.preferIPv4Stack=true")

# Run the JAR
options+=(-jar tracker-server.jar)

# Configuration file
options+=(/config/traccar.xml)

# Run the Traccar daemon
cd /opt/traccar || bashio::exit.nok "Failed to switch to Traccar directory"

exec java "${options[@]}"
